# templates/supervisord.conf.in
; Sample supervisor config file.

[unix_http_server]
file=${buildout:directory}/var/supervisor.sock   ; (the path to the socket file)

; uncomment to turn on HTTP access to control processes
;[inet_http_server]
;port = *:9001

[supervisord]
logfile=${buildout:directory}/var/log/supervisord.log 
logfile_maxbytes=50MB       
logfile_backups=10          
loglevel=info               
pidfile=${buildout:directory}/var/supervisord.pid 
nodaemon=false             
minfds=1024                
minprocs=200               


; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://${buildout:directory}/var/supervisor.sock ; use a unix:// URL  for a unix socket
history_file=~/.sc_history  ; use readline history if available

[program:node]
command = ${buildout:parts-directory}/nodejs/bin/yukeyukeyuke
redirect_stderr=true 
environment=PATH=${buildout:parts-directory}/nodejs/bin

[program:twisted]
command = ${buildout:bin-directory}/twistd -n -l - yuke
redirect_stderr=true

[program:docs-httpd]
command = ${buildout:bin-directory}/httpd -f ${buildout:directory}/conf/httpd.conf -DFOREGROUND
redirect_stderr=true
autostart = false

[program:docs]
command = python -m SimpleHTTPServer ${ports:public}
directory = ${buildout:directory}/www/html
redirect_stderr=true

[program:autodoc]
# maintains the 'nosier' inotify watching tool, to build the docs whenever the files change.
command = ${buildout:bin-directory}/nosier "make html"
redirect_stderr = true
directory = ${buildout:directory}/docs

[group:site]
programs=node,twisted,docs,autodoc
